# timers/update-wa-and-post/update-wa-and-post.timer
[Unit]
Description=Timer for update-wa-and-post.service

[Timer]
OnCalendar=*-*-* 14:00:00

[Install]
WantedBy=timers.target
