# timers/update-nz-and-post/update-nz-and-post.timer
[Unit]
Description=Timer for update-nz-and-post.service

[Timer]
OnCalendar=*-*-* 11:00:00

[Install]
WantedBy=timers.target
