# timers/update-vic-and-post/update-vic-and-post.timer
[Unit]
Description=Timer for update-vic-and-post.service

[Timer]
OnCalendar=*-*-* 08:30:00

[Install]
WantedBy=timers.target
