# timers/update-nsw-and-post/update-nsw-and-post.timer
[Unit]
Description=Timer for update-nsw-and-post.service

[Timer]
OnCalendar=*-*-* 11:00:00

[Install]
WantedBy=timers.target
