# timers/hello/hello.timer
[Unit]
Description=Timer for Hello

[Timer]
OnCalendar=*-*-* 21:38:00

[Install]
WantedBy=timers.target
