# timers/update-act-and-post/update-act-and-post.timer
[Unit]
Description=Timer for update-act-and-post.service

[Timer]
OnCalendar=*-*-* 11:45:00

[Install]
WantedBy=timers.target
